# systemd/walkman-mopidy.service
[Unit]
Description=Walkman: Mopidy (YouTube Music) audio server
After=network-online.target sound.target walkman-satellite.service
Wants=network-online.target walkman-satellite.service

[Service]
Type=simple
# Run as the brew user (reuses the verified config, YT auth, and audio access).
User=brew
Group=audio
Environment=HOME=/home/brew
# Shim dir holds sitecustomize.py which forces yt-dlp's android player client
# (avoids YouTube JS signature-solving -> fixes `next` and avoids a segfault).
Environment=PYTHONPATH=/home/brew/walkman/shim
ExecStart=/usr/bin/mopidy --config /home/brew/.config/mopidy/mopidy.conf
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
